# etc/services.d/telegraf/run
#!/usr/bin/with-contenv bash

set -eo pipefail

/usr/bin/telegraf \
    --config /etc/telegraf/telegraf.conf \
    2>&1 | awk '{print "[telegraf] " $0}'

