# etc/services.d/piaware2influx/run
#!/usr/bin/with-contenv bash

set -eo pipefail

/usr/bin/python3 \
    /piaware2influx.py \
    -ds $DUMP1090_HOST \
    2>&1 | awk '{print "[piaware2influx] " $0}'

